# scripts/podman/fluffbuzz.container.in
# FluffBuzz gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into fluffbuzz's ~/.config/containers/systemd/
# {{FLUFFBUZZ_HOME}} is replaced at install time.

[Unit]
Description=FluffBuzz gateway (rootless Podman)

[Container]
Image=fluffbuzz:local
ContainerName=fluffbuzz
UserNS=keep-id
Volume={{FLUFFBUZZ_HOME}}/.fluffbuzz:/home/node/.fluffbuzz
EnvironmentFile={{FLUFFBUZZ_HOME}}/.fluffbuzz/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
